# feelsjournal/config/feelsjournal.conf
[program:feelsjournal]
command = gunicorn server:app -b localhost:8000
directory = /home/www/feelsjournal
user = feels
